# scripts/rosbridge-websocket.service
[Unit]
Description=ROS rosbridge WebSocket (Foxglove Rosbridge, port 9090)
After=network-online.target
Wants=network-online.target
# roscore 须已由 start.sh / all.launch 拉起；脚本会轮询 ROS_MASTER_URI 直至可用

[Service]
Type=simple
User=orangepi
Group=orangepi
Environment=HOME=/home/orangepi
WorkingDirectory=/home/orangepi/newbot_ws
EnvironmentFile=-/etc/default/rosbridge
ExecStart=/home/orangepi/newbot_ws/scripts/run_rosbridge_websocket.sh
# roslaunch 若异常以 0 退出，on-failure 不会拉起重试
Restart=always
RestartSec=5
# 避免短时间内无限重启刷屏
StartLimitIntervalSec=300
StartLimitBurst=30

[Install]
WantedBy=multi-user.target
